# crontab
* * * * * cd /home/dns/dns god -c dns.god  -D >/dev/null 2>&1
